# supervisor/conf.d/piclaw.conf
; supervisor/conf.d/piclaw.conf – Supervisord program definition for piclaw.
;
; Runs piclaw via run-piclaw.sh with auto-restart. Logs to
; /var/log/piclaw/piclaw.stdout.log and .stderr.log.

[program:piclaw]
command=/usr/local/bin/run-piclaw.sh
user=agent
directory=/home/agent
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=INT
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/piclaw/piclaw.stdout.log
stderr_logfile=/var/log/piclaw/piclaw.stderr.log
environment=HOME="/home/agent",BUN_INSTALL="/usr/local/lib/bun",PATH="/workspace/.local/bin:/usr/local/lib/bun/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",PICLAW_WEB_TERMINAL_ENABLED="1"
