# supervisor/conf.d/piclaw.conf
[program:piclaw]
command=/usr/local/bin/run-piclaw.sh
user=agent
directory=/workspace
autostart=true
autorestart=true
startsecs=3
startretries=10
stopsignal=INT
stopwaitsecs=60
stdout_logfile=/var/log/piclaw/piclaw.stdout.log
stderr_logfile=/var/log/piclaw/piclaw.stderr.log
environment=HOME="/home/agent",BUN_INSTALL="/home/agent/.bun",PATH="/home/agent/.bun/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
